# doc/supervisor-conf.d/halva.conf
[program:halva]
directory=/home/pi/work/halva
command=node ./main.js
autostart=true
autorestart=true
stderr_logfile=/var/log/halva.err.log
stdout_logfile=/var/log/halva.out.log
user=pi
